# home/private_dot_config/autostart/protonvpn.desktop
[Desktop Entry]
Name=Proton VPN
Exec=flatpak run com.protonvpn.www
Type=Application
X-GNOME-Autostart-enabled=true